# debian/kdm.postrm
#!/bin/sh

set -e

case "$1" in
purge)
    update-rc.d kdm remove >/dev/null
    ;;

esac   

#DEBHELPER#
#KDEHELPER#
