# debian/postrm
#!/bin/sh

set -e

if [ "$1" = "purge" ]
then
  update-rc.d kdm remove >/dev/null
fi  

case "$1" in
remove) dpkg-divert --package kdebase --remove --rename \
	 --divert /etc/X11/Xsession.x11 /etc/X11/Xsession 
	chmod a-x /etc/menu-methods/kde 
	cp /etc/X11/window-managers /tmp/wm
	grep -v '\<kde\>' /tmp/wm > /etc/X11/window-managers
	;; 
esac   
