# debian/postinst
#!/bin/sh

set -e

if [ "$1" = "configure" ]
then
  update-rc.d kdm defaults 99 01 >/dev/null
fi  

chmod +x /etc/menu-methods/kde

window_manager_entry='/usr/X11R6/bin/kde'
if ! grep -q "^$window_manager_entry$" /etc/X11/window-managers 
then
	echo $window_manager_entry >> /etc/X11/window-managers 
        echo "Installed kde in the window-manager list."
fi
 
echo "Updating kdmrc ..."
(echo -n "s/WINDOWMANAGER/" ; \
cat /etc/X11/window-managers  |grep -v ^#|while read A ; \
do echo -n "`basename $A`;" ;  \
done ; \
echo "/") > /etc/kde/kdmrc.tmp
sed `cat /etc/kde/kdmrc.tmp` \
< /etc/kde/kdmrc.in > /etc/kde/kdmrc
rm /etc/kde/kdmrc.tmp

if test ! -L /var/spool/kdeapplnk; then
	ln -s /usr/share/applnk/debian_menu /var/spool/kdeapplnk
fi

#DEBHELPER#

# Starting kdm
/etc/init.d/kdm start
